# salt/packages/sonarqube/files/sonar.service
[Unit]
Description=SonarQube

[Service]
Type=forking
ExecStart=/opt/sonarqube-5.6.6/bin/linux-x86-64/sonar.sh start
ExecStop=/opt/sonarqube-5.6.6/bin/linux-x86-64/sonar.sh stop

[Install]
WantedBy=multi-user.target
